# sandboxd/sandboxd.service
[Unit]
Description=Tau OS Application Sandbox Daemon
Documentation=man:sandboxd(8)
After=network.target
Wants=network.target

[Service]
Type=simple
ExecStart=/usr/bin/sandboxd
Restart=on-failure
RestartSec=5
User=root
Group=root
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/tau/sandbox /tmp
PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictRealtime=true
RestrictSUIDSGID=true
LockPersonality=true
MemoryDenyWriteExecute=true

[Install]
WantedBy=multi-user.target 